# simple-keyboard-remapper.service
[Unit]
Description=Simple Keyboard Remapper
After=syslog.target

[Service]
Type=simple
ExecStart=/usr/bin/simple-keyboard-remapper /dev/input/event0
Restart=always

[Install]
WantedBy=multi-user.target
